# deploy/systemd/captainlabs-frontend.service.example
[Unit]
Description=CaptainLabs frontend
After=network.target captainlabs-api.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/srv/captainlabs/frontend
EnvironmentFile=/etc/captainlabs/frontend.env
ExecStart=/usr/bin/env npm run start
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/srv/captainlabs/frontend

[Install]
WantedBy=multi-user.target
